# meta-aspeed-sdk/meta-ast2600-pfr/recipes-support/spdm-emu/spdm-emu/bmc-attestation-emu.service
[Unit]
Description=BMC Attestation Emulator

[Service]
Restart=always
WorkingDirectory=/usr/share/spdm-emu
ExecStartPre=/usr/bin/slave-mqueue.sh 4 10
ExecStart=/usr/bin/spdm_responder_emu \
        --ver 1.0 \
        --cap CERT,CHAL,MEAS_SIG \
        --hash SHA_384 \
        --meas_spec DMTF \
        --meas_hash SHA_384 \
        --asym ECDSA_P384 \
        --req_asym ECDSA_P384 \
        --basic_mut_auth NO \
        --slot_id 0 \
        --mut_auth NO \
        --meas_sum NO \
        --meas_op ONE_BY_ONE \
        --meas_att HASH \
        --exe_mode SHUTDOWN \
        --exe_conn VER_ONLY,DIGEST,CERT,CHAL,MEAS \
        --trans MCTP
RestartSec=5
StartLimitInterval=0
Type=simple

[Install]
WantedBy=multi-user.target

