# meta-sophgo/recipes-phosphor/flash/phosphor-software-manager/sophgo-cpld-update@.service
[Unit]
Description=cpld update %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/sh -c '/usr/sbin/cpld-update.sh %i'
