# meta-sophgo/recipes-sophgo/sophgo-extend-file/files/set-fan-rate@.service
[Unit]
Description=Set fan speed to %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/sh -c '/usr/sbin/set-fan-rate.sh %i'
