# meta-sophgo-capricorn/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-host-bios@.service.in
[Unit]
Description=Flash Host Bios image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/sh -c '/usr/sbin/bios_update.sh /tmp/images/%i'
